# roles/rclone/templates/rclone-upload.path.j2
[Unit]
Description=Watch for new MoTeC log files to upload

[Path]
PathModified={{ log_buffer_path }}/gt7
PathModified={{ log_buffer_path }}/ams2

[Install]
WantedBy=multi-user.target
